# service/matter-netman/files/matter.init
#!/bin/sh /etc/rc.common

START=80

USE_PROCD=1
PROG=/usr/sbin/matter-network-manager-app

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_set_param user matter
	procd_set_param group matter
	procd_set_param stdout 1
	procd_set_param respawn
	procd_close_instance
}
